# crontab

# Deployment Monitor (every 5 minutes)
*/5 * * * * /root/.openclaw/trendwatcher/deploy-monitor.sh >> /root/.openclaw/trendwatcher/cron/logs/deploy-monitor.log 2>&1
